# service/cyracris.service
[Unit]
Description=CyraCRIS Flask (Cyra)
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/data/CyraCRIS
EnvironmentFile=/etc/cyracris.env
# Ajuste o caminho do Python do seu venv:
ExecStart=/data/CyraCRIS/venv/bin/gunicorn \
    --workers ${WORKERS:-4} \
    --bind 0.0.0.0:${PORT:-8000} \
    --timeout 60 \
    --access-logfile - \
    --error-logfile - \
    app:app
Restart=on-failure
RestartSec=5
# Hardening (opcional, pode afrouxar se precisar de escrita fora do diretório):
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
